# deploy/systemd/ironclaw-weixin-bridge.service
[Unit]
Description=IronClaw Weixin Bridge
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
WorkingDirectory=/root/ironclaw-weixin-bridge-repo
Environment=HOME=/root
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/node /root/ironclaw-weixin-bridge-repo/src/cli.mjs run
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
